# rs/etc/profile.d/00-fastfetch.sh
 [ -z "$PS1" ] && return

if ! [ $(id -u) = 0 ]; then if [ -f /usr/bin/fastfetch ]; then fastfetch; fi fi
