# rs/etc/update-motd.d/20-update-motd
#!/usr/bin/sh

echo "Hi there! You're connected to the machine" > /etc/motd
